# workers/whatsapp-bridge/com.techstartups.whatsapp-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.techstartups.whatsapp-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/YOUR_USERNAME/Projects/github/techstartups/distributed-infra/workers/whatsapp-bridge/bridge.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/Projects/github/techstartups/distributed-infra/workers/whatsapp-bridge</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>ORCHESTRATOR_URL</key>
    <string>http://100.97.176.37:8000</string>
    <key>INFRA_SECRET_KEY</key>
    <string>YOUR_SECRET_KEY_HERE</string>
    <key>WAHA_URL</key>
    <string>http://localhost:3000</string>
    <key>BRIDGE_PORT</key>
    <string>3001</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/whatsapp-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/whatsapp-bridge.err</string>
</dict>
</plist>
