# scripts/com.techstartups.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.techstartups.orchestrator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/michaelhansen/Projects/github/techstartups/distributed-infra/.venv/bin/uvicorn</string>
        <string>orchestrator.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/michaelhansen/Projects/github/techstartups/distributed-infra</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MACHINE_NAME</key>
        <string>michaels-macbook-pro</string>
        <key>MACHINE_ROLE</key>
        <string>orchestrator</string>
        <key>QUEUE_DB_PATH</key>
        <string>/Users/michaelhansen/Projects/github/techstartups/distributed-infra/data/queue.db</string>
        <key>SECRET_KEY</key>
        <string>e2de6168150d560393b14d1246725b0b9d8d6a453761dc7064bdbb3e6498726f</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/michaelhansen/Projects/github/techstartups/distributed-infra/logs/orchestrator.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/michaelhansen/Projects/github/techstartups/distributed-infra/logs/orchestrator.err</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
