# scripts/worker.service
# systemd user-level unit for ThinkPad (Ubuntu) worker — no sudo needed.
# Install: cp scripts/worker.service ~/.config/systemd/user/infra-worker.service
#          systemctl --user daemon-reload
#          systemctl --user enable --now infra-worker
#          loginctl enable-linger $USER   # survive logout

[Unit]
Description=Distributed Infra Worker
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/michaelhansen/Projects/distributed-infra
EnvironmentFile=/home/michaelhansen/Projects/distributed-infra/.env
Environment=LD_LIBRARY_PATH=/home/michaelhansen/lib
ExecStart=/home/michaelhansen/Projects/distributed-infra/.venv/bin/uvicorn worker.main:app --host 0.0.0.0 --port 8001
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
