# scripts/worker.service
# systemd unit for ThinkPad (Ubuntu) worker
# Install: sudo cp scripts/worker.service /etc/systemd/system/infra-worker.service
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now infra-worker

[Unit]
Description=Distributed Infra Worker
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/distributed-infra
EnvironmentFile=/home/%i/distributed-infra/.env
ExecStart=/home/%i/distributed-infra/.venv/bin/uvicorn worker.main:app --host 0.0.0.0 --port 8001
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
